# ansible/roles/watcher/files/watcher.service
[Unit]
Description=Simple File Watcher

[Service]
Type=oneshot
ExecStart=/bin/systemctl restart locker
