# ansible/roles/watcher/files/watcher.path
[Unit]
Description=Monitor lockfile file for changes

[Path]
PathChanged=/locks/lockfile.lock

[Install]
WantedBy=multi-user.target
